# systemd/golden-plate.service
[Unit]
Description=Golden Plate Kiosk (Xorg + Openbox + kiosk-controller)
Wants=network-online.target
After=network-online.target systemd-user-sessions.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/golden-plate

# optional (wird an X/Child-Prozesse weitergereicht)
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/user/.Xauthority
EnvironmentFile=-/etc/default/kiosk-controller

# VT1 sicher claimen, dann X-Session im User-Kontext starten
ExecStart=/usr/bin/openvt -f -c 1 -- /bin/su - user -c /usr/bin/startx\ /opt/golden-plate/scripts/kiosk-session.sh\ --\ :0\ -nolisten\ tcp\ vt1

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target