# systemd/kiosk.service
[Unit]
Description=Golden Plate - Kiosk (X11 + Chromium)
After=kiosk-web.service
Wants=kiosk-web.service
StartLimitIntervalSec=0

[Service]
User=kiosk
WorkingDirectory=/home/kiosk

Environment=HOME=/home/kiosk
Environment=DISPLAY=:0

ExecStart=/usr/bin/startx /opt/golden-plate/bin/session.sh -- :0 -nolisten tcp vt7

# startx/xinit endet beim Stop oft mit 1 oder 143 (SIGTERM).
SuccessExitStatus=0 1 143
# optional: SIGKILL nach Timeout explizit erlauben (Default ist meist yes, aber klarer)
SendSIGKILL=yes

Restart=on-failure
RestartSec=2

TimeoutStartSec=30
TimeoutStopSec=10
KillMode=control-group

[Install]
WantedBy=multi-user.target