# systemd/kiosk-monitor.timer
[Unit]
Description=Golden Plate - Run kiosk health monitor periodically

[Timer]
OnBootSec=60
OnUnitActiveSec=900
AccuracySec=5
Unit=kiosk-monitor.service

[Install]
WantedBy=timers.target