# systemd/kiosk-monitor.service
[Unit]
Description=Golden Plate - Kiosk Health Monitor (X11/Chromium/LocalWeb)
After=kiosk.service kiosk-web.service

[Service]
Type=oneshot
ExecStart=/opt/golden-plate/bin/kiosk-monitor.sh